# debian-src/trusty-quilt/qigong-1.10.10/debian/qigong-nomc.postrm
#!/bin/bash

set -e

if [ "$1" = purge ]; then
    update-rc.d qigong remove
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
