# qicollect.rc.proto
#!/bin/sh


### BEGIN INIT INFO
# Provides:          qicollect
# Required-Start:    $syslog $time $network $local_fs
# Required-Stop:     $syslog $time $network $local_fs
# Should-Start:	     qigong
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: status collecting daemon
# Description:       manages the status collecting daemon that fetches infos from qigong(s)
### END INIT INFO


PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DESC="status collecting daemon"
NAME="qicollect"
DAEMON="@@PREFIX@@/sbin/${NAME}"
PIDFILE="/var/run/${NAME}.pid"
CONFFILE="/etc/qicollect.conf"
SCRIPTNAME="/etc/init.d/${NAME}"
USER=qicollect

d_start() {
	: > "${PIDFILE}"
	chown "${USER}" "${PIDFILE}"
	su - "${USER}" -c "${DAEMON} -pidfile \"${PIDFILE}\""
}

d_stop() {
	if [ ! -e "${PIDFILE}" ]
	then
	    echo
	    echo " no pid file found (${PIDFILE})"
	    exit 1
	fi
	PID=`cat "${PIDFILE}"`
	EXENAME=`ps -p "${PID}" | grep "${PID}" | cut -b25-`
	if [ E"${EXENAME}" = E"${NAME}" ]
	then
	    kill "${PID}"
	    if [ -e "${PIDFILE}" ]
	    then
		echo " cleaning ${PIDFILE} "
		rm "${PIDFILE}"
	    fi
	else
	    echo
	    echo " pid file doesn't match daemon"
	    exit 1
	fi
}

d_gracefulrestart() {
	if [ ! -e "${PIDFILE}" ]
	then
	    echo
	    echo " no pid file found (${PIDFILE})"
	    exit 1
	fi
	PID=`cat "${PIDFILE}"`
	EXENAME=`ps -p "${PID}" | grep "${PID}" | cut -b25-`
	if [ E"${EXENAME}" = E"${NAME}" ]
	then
	    kill -s USR1 "${PID}"
	else
	    echo
	    echo " pid file doesn't match daemon"
	    exit 1
	fi
	
}

d_status() {
	if [ ! -e "${PIDFILE}" ]
	then
	    echo "${NAME} is not running"
	    exit 3
	fi
	PID=`cat "${PIDFILE}"`
	EXENAME=`ps -p "${PID}" | grep "${PID}" | cut -b25-`
	if [ E"${EXENAME}" = E"${NAME}" ]
	then
	    echo "${NAME} is running"
	    exit 0
	else
	    echo "${NAME} is not and pidfile ${PIDFILE} exists"
	    exit 1
	fi
}

# Gracefully exit if the package has been removed.
test -x "${DAEMON}" || exit 0
# Gracefully exit if the package isn't configured
if [ ! -e "${CONFFILE}" ]
then
    echo "skipping $DESC: $NAME, no conf-file installed yet."
    exit 0
fi

case "$1" in
  start)
	echo -n "Starting $DESC: $NAME"
	d_start
	echo "."
	;;
  stop)
	echo -n "Stopping $DESC: $NAME"
	d_stop
	echo "."
	;;
  reload|force-reload)
	echo -n "Reloading $DESC configuration via signal sent"
	d_gracefulrestart
	echo "done."
  	;;
  restart)
	echo -n "Restarting $DESC: $NAME"
	d_stop
	sleep 1
	d_start
	echo "."
	;;
  status)
	d_status && exit 0 || exit $?
	;;
  *)
	echo "Usage: $SCRIPTNAME {reload|start|stop|restart|force-reload|status}" >&2
	exit 1
	;;
esac

exit 0

