# qigong.rc.proto
#!/bin/sh

### BEGIN INIT INFO
# Provides:          qigong
# Required-Start:    $syslog $time $network $local_fs $remote_fs
# Required-Stop:     $syslog $time $network $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: local status collecting daemon
# Description:       manages the local status collecting daemon that spread infos for qicollect(s)
### END INIT INFO


PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DESC="status reporting daemon"
NAME="qigong"
DAEMON="@@PREFIX@@/sbin/${NAME}"
PIDFILE="/var/run/${NAME}.pid"
SCRIPTNAME="/etc/init.d/${NAME}"

if [ ! -f /etc/default/"${NAME}" ]
then
    echo "$@" " : missing /etc/default/${NAME}"
    exit 1
fi
source /etc/default/"${NAME}"

d_start() {
	"${DAEMON}" -pidfile "${PIDFILE}" -localkey="${KEYFILE}"
}

d_stop() {
	if [ ! -e "${PIDFILE}" ]
	then
	    echo
	    echo " no pid file found (${PIDFILE})"
	    return 1
	fi
	PID=`cat "${PIDFILE}"`
	EXENAME=`ps -p "${PID}" | grep "${PID}" | cut -b25-`
	if [ E"${EXENAME}" = E"${NAME}" ]
	then
	    kill "${PID}"
	    if [ -e "${PIDFILE}" ]
	    then
		echo " cleaning ${PIDFILE} "
		rm "${PIDFILE}"
	    fi
	else
	    echo
	    echo " pid file doesn't match daemon"
	    return 1
	fi
}

d_status() {
	if [ ! -e "${PIDFILE}" ]
	then
	    echo "${NAME} is not running"
	    return 3
	fi
	PID=`cat "${PIDFILE}"`
	EXENAME=`ps -p "${PID}" | grep "${PID}" | cut -b25-`
	if [ E"${EXENAME}" = E"${NAME}" ]
	then
	    echo "${NAME} is running"
	    return 0
	else
	    echo "${NAME} is not and pidfile ${PIDFILE} exists"
	    return 1
	fi
}

# Gracefully exit if the package has been removed.
test -x "${DAEMON}" || exit 0

case "$1" in
  start)
	echo -n "Starting $DESC: $NAME"
	d_start
	echo "."
	;;
  stop)
	echo -n "Stopping $DESC: $NAME"
	d_stop
	echo "."
	;;
    #	  reload)
    #		echo -n "Reloading $DESC configuration via signal sent"
    #		d_gracefulrestart
    #		echo "done."
    #	  	;;
  restart|force-reload)
	#
	#	If the "reload" option is implemented, move the "force-reload"
	#	option to the "reload" entry above. If not, "force-reload" is
	#	just the same as "restart".
	#
	echo -n "Restarting $DESC: $NAME"
	d_stop
	sleep 1
	d_start
	echo "."
	;;
  status)
	d_status && exit 0 || exit $?
	;;
  *)
	echo "Usage: $SCRIPTNAME {start|stop|restart|force-reload|status}" >&2
	exit 1
	;;
esac

exit 0

