# DEBIAN-qigong-nomc/prerm
#!/bin/bash

set -e

if which invoke-rc.d >/dev/null 2>&1; then
    invoke-rc.d qigong stop
else
    /etc/init.d/qigong stop
fi

